# deploy/backup/pg_backup.service
[Unit]
Description=windy-registry Postgres backup → S3
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/windy-registry-backup.env
ExecStart=/opt/windy-registry/deploy/backup/pg_backup.sh
StandardOutput=journal
StandardError=journal
